# yocto-meta-layers/meta-bare-metal-router/recipes-core/startup-systemd/files/start-routershell.service
[Unit]
Description=Routershell Python CLI Service
After=getty.target

[Service]
Type=simple
ExecStart=/etc/routershell/start.sh
Restart=always
StandardInput=tty-force

[Install]
WantedBy=multi-user.target
